# backend/tests/fixtures/scheduled_tasks/appdata_script_task.xml
<?xml version="1.0" encoding="UTF-8"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2026-05-03T12:00:00Z</Date>
    <Author>alex</Author>
    <URI>\Chrome Update Service</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="Author">
      <UserId>alex</UserId>
    </Principal>
  </Principals>
  <Settings>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
  </Settings>
  <Triggers>
    <BootTrigger>
      <Enabled>true</Enabled>
    </BootTrigger>
  </Triggers>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-File C:\Users\alex\AppData\Local\Temp\runme.ps1</Arguments>
    </Exec>
  </Actions>
</Task>
