# backend/tests/fixtures/scheduled_tasks/suspicious_powershell_encoded.xml
<?xml version="1.0" encoding="UTF-8"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2026-05-03T11:30:00Z</Date>
    <Author>alex</Author>
    <URI>\Windows Update Monitor</URI>
  </RegistrationInfo>
  <Principals>
    <Principal id="Author">
      <UserId>alex</UserId>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <Enabled>true</Enabled>
    <Hidden>true</Hidden>
  </Settings>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
    </LogonTrigger>
  </Triggers>
  <Actions Context="Author">
    <Exec>
      <Command>C:\Windows\System32\WindowsPowerShell\v1.0\powershell.exe</Command>
      <Arguments>-NoP -EncodedCommand SQBFAFgA</Arguments>
      <WorkingDirectory>C:\Users\alex\AppData\Local\Temp</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
